# contrib/systemd/irom-tpu-scheduler.service
[Unit]
Description=Personal IROM TPU queue scheduler
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
Environment=PATH=%h/google-cloud-sdk/bin:%h/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=TPU_QUEUE_USER=%u
WorkingDirectory=%h
ExecStart=%h/.local/bin/tpu scheduler --focus-user=%u --scan-interval=30
Restart=always
RestartSec=10
KillMode=control-group
TimeoutStopSec=30

[Install]
WantedBy=default.target
